# daemons/dnssec/ipa-ods-exporter.socket
[Socket]
ListenStream=/var/run/opendnssec/engine.sock

[Install]
WantedBy=sockets.target
