# daemons/dnssec/ipa-ods-exporter.socket.in
[Socket]
ListenStream=@runstatedir@/opendnssec/engine.sock
SocketUser=@ODS_USER@
SocketGroup=@ODS_GROUP@

[Install]
WantedBy=sockets.target
