# daemons/dnssec/ipa-ods-exporter.socket.in
[Socket]
ListenStream=@localstatedir@/run/opendnssec/engine.sock

[Install]
WantedBy=sockets.target
