# daemons/dnssec/ipa-ods-exporter.socket.in
[Socket]
ListenStream=@runstatedir@/opendnssec/engine.sock

[Install]
WantedBy=sockets.target
