# daemons/ipa-otpd/ipa-otpd.socket.in
[Unit]
Description=ipa-otpd socket

[Socket]
ListenStream=@krb5rundir@/DEFAULT.socket
ExecStopPre=-@UNLINK@ @krb5rundir@/DEFAULT.socket
SocketMode=0600
Accept=true

[Install]
WantedBy=krb5kdc.service
