# init/systemd/httpd.service
.include /usr/lib/systemd/system/httpd.service

[Service]
Environment=KRB5CCNAME=/var/run/httpd/krbcache/krb5ccache
