# init/systemd/ipa-custodia.service.in
[Unit]
Description=IPA Custodia Service

[Service]
Type=notify
ExecStart=@libexecdir@/ipa/ipa-custodia @IPA_SYSCONF_DIR@/custodia/custodia.conf
PrivateTmp=yes
Restart=on-failure
RestartSec=60s

[Install]
WantedBy=multi-user.target
