# init/systemd/ipa-ccache-sweep.timer.in
[Unit]
Description=Remove Expired Kerberos Credential Caches

[Timer]
OnUnitActiveSec=12h

[Install]
WantedBy=timers.target
