# init/systemd/ipa.service.in
[Unit]
Description=Identity, Policy, Audit
Wants=network.target
Wants=gssproxy.service
After=network.target

[Service]
Type=oneshot
ExecStart=@sbindir@/ipactl start
ExecStop=@sbindir@/ipactl stop
RemainAfterExit=yes
TimeoutSec=0

[Install]
WantedBy=multi-user.target
