# client/systemd/ipa-epn.service.in
[Unit]
Description=Execute IPA Expiring Password Notification (EPN)

[Service]
Type=simple
ExecStart=@sbindir@/ipa-epn

[Install]
WantedBy=multi-user.target
